# os-template/etc/profile.d/mock.sh
# This causes any mock builds to go to a directory ending in the current username
alias mock="mock --uniqueext=$USER"

